# src/system/monitor-serial.service
# Copy from: src/system/monitor-serial.service
# Install to: /etc/systemd/system/monitor-serial.service
#
# Behavior: Restarts only on crash (not on manual stop).
# If it crashes 3 times within 30 seconds, systemd gives up.
#
# To install and enable:
# 1. sudo cp src/system/monitor-serial.service /etc/systemd/system/
# 2. sudo systemctl daemon-reload
# 3. sudo systemctl enable monitor-serial.service
# 4. sudo systemctl start monitor-serial.service
# To stop or restart:
# 5. sudo systemctl stop monitor-serial.service
# 6. sudo systemctl restart monitor-serial.service

[Unit]
Description=LED Table Serial Monitor
After=network.target

[Service]
ExecStart=/home/tcledtable/projects/py_tc_led_table/.venv/bin/python scripts/monitor_serial.py
WorkingDirectory=/home/tcledtable/projects/py_tc_led_table
User=tcledtable
Group=tcledtable
Environment=PYTHONUNBUFFERED=1
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=on-failure
RestartSec=1
StartLimitBurst=3
StartLimitIntervalSec=30

[Install]
WantedBy=multi-user.target
